# scripts/com.claw.guitar-platform.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claw.guitar-platform.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/claw/GuitarPlatform_v2/scripts/run_backend_service.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/claw/GuitarPlatform_v2</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/claw/Library/Logs/GuitarPlatform/backend.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/claw/Library/Logs/GuitarPlatform/backend.error.log</string>
</dict>
</plist>
